# pc_installation/vpn_connection.service
[Unit]
After=multi-user.target

[Service]
ExecStart=/usr/local/bin/vpn_connection.sh

[Install]
WantedBy=multi-user.target
